# root/hidpi/etc/profile.d/XRes.sh
#!/bin/sh
xrdb /etc/Xresources
export XCURSOR_SIZE=36
